# deploy/systemd/entrada-pro-audit.service
[Unit]
Description=ENTRADA-PRO Auditor (gera audit.json + gate)
After=network.target

[Service]
Type=oneshot
User=roteiro_ds
WorkingDirectory=/home/roteiro_ds/ENTRADA-PRO
Environment=DATA_DIR=/home/roteiro_ds/AUTOTRADER-PRO/data
ExecStart=/bin/bash /home/roteiro_ds/ENTRADA-PRO/scripts/run_audit.sh
ExecStartPost=/bin/bash /home/roteiro_ds/ENTRADA-PRO/scripts/audit_gate.sh http://127.0.0.1:8090
