# deploy/systemd/entrada-pro-audit.timer
[Unit]
Description=ENTRADA-PRO Auditor Gate timer (cada 5 min)

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Unit=entrada-pro-audit.service

[Install]
WantedBy=timers.target
