# deploy/systemd/entrada-pro-worker.timer
[Unit]
Description=ENTRADA-PRO Worker timer (5 min)

[Timer]
OnBootSec=30s
OnUnitActiveSec=5min
Unit=entrada-pro-worker.service

[Install]
WantedBy=timers.target
