# deploy/systemd/entrada-pro-worker.timer
[Unit]
Description=ENTRADA-PRO Worker timer (5 min)

[Timer]
OnBootSec=60
OnUnitActiveSec=300
Unit=entrada-pro-worker.service

[Install]
WantedBy=timers.target
