# deploy/systemd/entrada-pro-audit.service
[Unit]
Description=ENTRADA-PRO Auditor Gate (bloqueia se falhar)
After=network.target

[Service]
Type=oneshot
User=roteiro_ds
WorkingDirectory=/home/roteiro_ds/ENTRADA-PRO
ExecStart=/bin/bash /home/roteiro_ds/ENTRADA-PRO/scripts/audit_gate.sh http://127.0.0.1:8090
